# packaging/systemd/snooze-server.service
[Unit]
Description=Snooze main server daemon
Documentation=https://github.com/snoozeweb/snooze
After=network-online.target
Wants=network-online.target
# Snooze supports MongoDB, PostgreSQL and SQLite backends. The packaging
# does not pin a specific database service here so any of the three may
# be used. If your deployment uses a co-located DB, add an override:
#   systemctl edit snooze-server
# and set e.g. `Wants=mongodb.service` or `Wants=postgresql.service`.

[Service]
Type=simple
User=snooze
Group=snooze
# The default SQLite/file backend writes a relative ./db.sqlite; root the CWD at
# the writable state dir (also a ReadWritePath below) so the file lands there
# rather than in the read-only / under ProtectSystem=strict.
WorkingDirectory=/var/lib/snooze
# -config takes a DIRECTORY of section YAML files (core.yaml, general.yaml, …),
# not a single file. An empty dir is fine: the server boots on defaults + env
# (the default SQLite DB then lives at /var/lib/snooze/db.sqlite).
ExecStart=/usr/bin/snooze-server -config /etc/snooze/server
Restart=on-failure
RestartSec=5s

# Hardening
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/snooze /var/log/snooze
# /run/snooze hosts the admin Unix socket (snoozectl, etc.).
RuntimeDirectory=snooze
RuntimeDirectoryMode=0750

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
